# fwatchd.service
[Unit]
Description=fwatchd - a file watching daemon

[Service]
Type=forking
ExecStart=/usr/sbin/fwatchd

[Install]
WantedBy=default.target
